# Vegtrug_IoTHub/daemonIoTHub.service
[Unit]
Description="Mi Flora, Vegtrug plant monitor"
After=multi-user.target

[Service]
Type=simple
ExecStart=/home/pi/dev/Vegtrug/Vegtrug_IoTHub/MiFloraMain.py
User=pi
WorkingDirectory=/home/pi/dev/Vegtrug/Vegtrug_IoTHub
Restart=on-failure

[Install]
WantedBy=multi-user.target